# Maquinas/A01-Perdida_de_Control_de_Acceso/multi_stage_docker_DEMO/Cron.service
# /etc/systemd/system/cron.service
[Unit]
Description=Cron Service
After=network.target

[Service]
ExecStart=/usr/sbin/cron -f
Restart=always

[Install]
WantedBy=multi-user.target
